# supervisord.conf
[supervisord]
user=root
nodaemon=true
pidfile=/dev/null
logfile=/dev/null
logfile_maxbytes=0

[program:example]
command=/bin/sleep 10
autorestart=true
